# contribute/pg-backup.timer
[Unit]
Description=Run PostgreSQL backup daily

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
# Or for a specific time: OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target