# contribute/pg-backup.service
[Unit]
Description=PostgreSQL 17 Incremental Backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Path to your backups on the host
# Ensure this directory exists and has correct permissions
ExecStartPre=-/usr/bin/mkdir -p /opt/pg_backups
# Path to your environment file
EnvironmentFile=/etc/pg-backup.env
# Pull the image if it's not present
ExecStartPre=/usr/bin/podman pull docker.io/library/pg17-backup:latest
# Run the backup container
ExecStart=/usr/bin/podman run --rm \
  --env-file=/etc/pg-backup.env \
  -v /opt/pg_backups:/backups \
  --name pg17-backup-job \
  docker.io/library/pg17-backup:latest